# profile.d/my-bin.sh
export PATH=/home/osily/.bin:/usr/lib/colorgcc/bin:$PATH
